# supervisorctl/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
logfile=/var/log/ray/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/ray

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# Gunicorn — Django (config.wsgi). Ensure /opt/ray/run exists and is writable by the app user.
# Loads /opt/ray/.env (see deploy/production.env.example) so DJANGO_SECRET_KEY etc. apply.
# Gunicorn binary: use /opt/ray/.venv/bin/gunicorn if your venv is named .venv instead of venv.
[program:gunicorn]
command=/bin/bash -c 'set -a; [ -f /opt/ray/.env ] && . /opt/ray/.env; set +a; exec /opt/ray/venv/bin/gunicorn config.wsgi:application --bind unix:/opt/ray/run/django.sock --umask 007 --workers 2 --threads 2 --log-level info --timeout 120 --max-requests 1000 --max-requests-jitter 50 --graceful-timeout 30 --keep-alive 5 --worker-class gthread --worker-connections 1000'
directory=/opt/ray
user=deploy
group=www-data
autostart=true
autorestart=true
stderr_logfile=/var/log/ray/gunicorn.err.log
stdout_logfile=/var/log/ray/gunicorn.out.log
environment=PYTHONPATH="/opt/ray",DJANGO_SETTINGS_MODULE="config.settings"
startsecs=5
startretries=3

[group:ray]
programs=gunicorn
